# home/bilabila/.config/systemd/user/arknights.service
[Unit]
Description=restart idehelper

[Service]
Type=simple
KillMode=process
ExecStart=/home/bilabila/bin/scripts/restart_idehelper
